# scripts/com.arroyo.blog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Daily sermon-blog generator. Runs every day at 9:15am and generates a draft for
     the LATEST sermon as soon as YouTube has captioned it (idempotent: one draft per
     sermon). If the Mac is asleep/off at 9:15, launchd runs the missed job on next wake.
     Install:
       cp scripts/com.arroyo.blog.plist ~/Library/LaunchAgents/
       launchctl load ~/Library/LaunchAgents/com.arroyo.blog.plist -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.arroyo.blog</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/dakotayates/ai-os/ventures/arroyo-church-redesign/scripts/run-blog.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key><integer>9</integer>
    <key>Minute</key><integer>15</integer>
  </dict>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/dakotayates/Library/Logs/arroyo-blog.launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/dakotayates/Library/Logs/arroyo-blog.launchd.log</string>
</dict>
</plist>
